# systemd/etc/systemd/user/window-manager.service.d/00-custom.conf
[Unit]
Requisite=desktop.target
;After=desktop.target
;Requires=xorg.service xorg.socket
;After=xorg.service xorg.socket
;Requires=desktop.target

# vim: set ft=systemd :
