# systemd/etc/systemd/user/window-manager.service.d/00-desktop.target.conf
[Unit]
Requisite=desktop.target
After=desktop.target

# vim: set ft=systemd :
