# systemd/etc/systemd/user/window-manager.service.d/20-custom.conf
[Unit]
Requisite=desktop.target
After=desktop.target

Wants=compositor.service notification.service 
Before=compositor.service notification.service

Wants=wallpaper.service

# vim: set ft=systemd :
