# systemd/etc/systemd/user/chromium.service
[Unit]
Description=Chromium
Requires=chromium.socket
After=desktop.target chromium.socket

[Service]
Type=simple
ExecStart=/usr/bin/chromium
