# systemd/etc/systemd/user/desktop.target
[Unit]
Description=Desktop
After=xorg.target window-manager.service
Requires=dbus.socket xorg.target window-manager.service

[Install]
Alias=default.target
