# systemd/etc/systemd/user/desktop.target
[Unit]
Description=Desktop
Requires=graphical.target
Wants=window-manager.service wallpaper.service
After=graphical.target
#After=window-manager.service
