# systemd/etc/systemd/user/desktop.target
[Unit]
Description=Desktop Interface
AllowIsolate=no


# NOTE: socket activation is broken in xorg 1.18, so I need to require
# xorg.service as well.
Requires=basic.target dbus.socket xorg.service xorg.socket
After=basic.target dbus.socket xorg.service xorg.socket

Requires=window-manager.service
Before=window-manager.service

Wants=xinit.service
Wants=compositor.service
Wants=notification.service
Wants=wallpaper.service
