# systemd/etc/systemd/user/desktop.target
[Unit]
Description=Desktop Interface
AllowIsolate=no

PartOf=window-manager.service xorg.socket xorg.service
BindsTo=window-manager.service

Requires=basic.target dbus.socket window-manager.service
Before=window-manager.service

Wants=xinit.service
Wants=compositor.service
Wants=notification.service
Wants=wallpaper.service
