# systemd/etc/systemd/user/gpg-agent.service
[Unit]
Description=GPG Agent service

[Service]
Type=forking
ExecStart=/usr/bin/gpg-agent --verbose --daemon

[Install]
WantedBy=basic.target
