# systemd/etc/systemd/user/i3wm.service
[Unit]
Description=i3 tiling window manager

Requires=xorg.service
After=xorg.service

BindsTo=xorg.service xorg.socket
PartOf=xorg.service xorg.socket

[Service]
Environment=DISPLAY=:0
Environment=DE=i3
ExecStart=/usr/bin/i3
ExecReload=/usr/bin/i3-msg reload

[Install]
Alias=window-manager.service
