# systemd/etc/systemd/user/pidgin.service
[Unit]
Description=Pidgin
Requisite=desktop.target

[Service]
ExecStart=/usr/bin/pidgin

[Install]
WantedBy=desktop.target
