# systemd/etc/systemd/user/wallpaper.timer
[Unit]
Description=Random wallpaper timer
Requisite=desktop.target

[Timer]
OnUnitActiveSec=1800s
